# deploy/muscledb-web.service
[Unit]
Description=muscledb-web (anatomy reference API + web UI)
After=network.target

[Service]
Type=simple
# Runs as your login user so you can scp files into place without sudo/chown.
User=bp
WorkingDirectory=/opt/muscledb-web
ExecStart=/usr/bin/node server.js
Environment=NODE_ENV=production
Environment=PORT=5000
Environment=HOST=127.0.0.1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
